# usr.bin/uucp/uudemon.crontab
#
#ident  "@(#)uudemon.crontab 1.1 92/07/30"
#
48 8,12,16 * * * /etc/uucp/uudemon.admin
45 23 * * * /etc/uucp/uudemon.cleanup
0 * * * * /etc/uucp/uudemon.poll
11,41 * * * * /etc/uucp/uudemon.hour
